# ansible/templates/supervisor-backend.conf.j2
[program:cybersecurity-{{ environment }}-backend]
command={{ env_dirs[environment].backend }}/venv/bin/python app.py
directory={{ env_dirs[environment].backend }}
user={{ app_user }}
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile={{ env_dirs[environment].logs }}/backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=FLASK_ENV="{{ environment }}",DATABASE_URL="postgresql://pentest_{{ environment }}_user:{{ db_password }}@localhost:{{ postgres_port }}/pentest_{{ environment }}"

